# etc/supervisor/conf.d/quasselcore.conf
[program:quasselcore]
command=quasselcore --configdir=/var/lib/quassel

